# Services/.config/systemd/user/setWallpaper.service
[Unit]
Description=Set Wallpaper Periodically

[Service]
Type=oneshot
ExecStart=sh -c "dash $HOME/.local/share/bin/setWallpaper.sh"

