# Services/.config/systemd/user/setWallpaper.timer
[Unit]
Description=Run Script at sunset

[Timer]
OnCalendar=*-*-* *:*:00
OnActiveSec=3s
Unit=setWallpaper.service

[Install]
WantedBy=timers.target

